# LocalDev/launchagents/com.flownote.local.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>KeepAlive</key>
  <true/>
  <key>Label</key>
  <string>com.flownote.local.api</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/truds/Projects/Project/FlowNote/LocalDev/start-local-service.sh</string>
    <string>api</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/truds/Projects/Project/FlowNote/services/api</string>
  <key>StandardOutPath</key>
  <string>/tmp/flownote-api.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/flownote-api.error.log</string>
</dict>
</plist>
